# crbuild.service
[Unit]
Description=Custom Rom Build Script

[Service]
WorkingDirectory=%h/los
ExecStart=%h/bin/crbuilder

[Install]
WantedBy=default.target
